# internal/template/templates/gdrivesync/com.dotfiles.gdrive-sync.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{.Label}}</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{.DotfilesPath}}</string>
    <string>gdrive-sync</string>
    <string>{{.Action}}</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd inherits a minimal PATH that finds Apple's old /usr/bin/rsync
         (openrsync, "rsync 2.6.9 compatible") instead of Homebrew's modern
         rsync 3.x. List Homebrew prefixes first so `--info=stats2,progress2`
         and other 3.x-only flags resolve correctly on both arches. -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key>
  <integer>{{.Interval}}</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{.LogFile}}</string>
  <key>StandardErrorPath</key>
  <string>{{.LogFile}}</string>
</dict>
</plist>
