# internal/template/templates/gdrivesync/dotfiles-gdrive-sync.service.tmpl
[Unit]
Description={{.Description}}

[Service]
Type=oneshot
ExecStart={{.DotfilesPath}} gdrive-sync {{.Action}}
