# internal/template/templates/gdrivesync/dotfiles-gdrive-sync.timer.tmpl
[Unit]
Description={{.Description}} (timer)

[Timer]
OnBootSec=60
OnUnitActiveSec={{.Interval}}s
Unit={{.ServiceName}}

[Install]
WantedBy=timers.target
