# internal/template/templates/gdrivesync/dotfiles-gdrive-sync.service.tmpl
[Unit]
Description=Periodic gdrive-sync (workspace ⟷ gdrive-workspace mirror)

[Service]
Type=oneshot
ExecStart={{.DotfilesPath}} gdrive-sync sync
