# internal/template/templates/sync/rclone-bisync.service.tmpl
[Unit]
Description=Sync workspace with Google Drive

[Service]
Type=oneshot
ExecStart={{.DotfilesPath}} sync
