# internal/template/templates/sync/rclone-bisync.timer.tmpl
[Unit]
Description=Periodic rclone bisync

[Timer]
OnBootSec=60
OnUnitActiveSec={{.Interval}}s
Unit=rclone-bisync.service

[Install]
WantedBy=timers.target
