# internal/template/templates/sync/rclone-bisync.service.tmpl
[Unit]
Description=rclone bisync workspace to Google Drive

[Service]
Type=oneshot
ExecStart={{.RclonePath}} bisync {{.LocalPath}} {{.RemotePath}} \
  --filter-from {{.FilterFile}} \
  --conflict-resolve newer \
  --conflict-loser num \
  --resilient \
  --recover \
  --max-lock 15m \
  --tpslimit 20 \
  --retries 5 \
  --fast-list \
  --drive-skip-dangling-shortcuts \
  --log-file {{.LogFile}} \
  --log-level INFO
