# internal/template/templates/gsync/dotfiles-gdrive-sync.service.tmpl
[Unit]
Description={{.Description}}

[Service]
Type=oneshot
ExecStart={{.DotfilesPath}} gsync {{.Action}} --mode={{.Mode}}
