# internal/template/templates/rsync/com.dotfiles.workspace-sync.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dotfiles.workspace-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{.DotfilesPath}}</string>
    <string>sync</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd's default PATH finds Apple's old /usr/bin/rsync (openrsync,
         "rsync 2.6.9 compatible"). Put Homebrew first so modern rsync 3.x
         resolves on both Apple Silicon and Intel Macs. -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key>
  <integer>{{.Interval}}</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{.LogFile}}</string>
  <key>StandardErrorPath</key>
  <string>{{.LogFile}}</string>
</dict>
</plist>
