# internal/template/templates/rsync/dotfiles-sync.service.tmpl
[Unit]
Description=Sync workspace binaries with remote server

[Service]
Type=oneshot
ExecStart={{.DotfilesPath}} sync
